# systemd/uwsgi.service.d/posthorn.conf
[Unit]
After=posthorn-database-setup.service
